# dev/zouip_sender.service
[Unit]
Description=zouip

[Service]
ExecStart=/home/tonton/Taf/zouip/dev/clipboard_send.sh

[Install]
WantedBy=default.target
