# dev/zouip_sender.service
[Unit]
Description=zouip

[Service]
ExecStart=/home/tonton/Taf/zouip/clipboard_receive.sh

[Install]
WantedBy=default.target
