# dev/python/resources/zouip_sender.service
[Unit]
Description=zouip

[Service]
ExecStart=python3

[Install]
WantedBy=default.target
